# scripts/systemd/shiori-inbox-refresh.timer
[Unit]
Description=Run Shiori inbox refresh every 2 hours

[Timer]
OnBootSec=5m
OnUnitActiveSec=2h
Unit=shiori-inbox-refresh.service
Persistent=true

[Install]
WantedBy=timers.target
